# target/linux/ralink/wp688/wp688_MTK4110_130426/base-files/etc/init.d/channellist
#!/bin/sh /etc/rc.common

START=20
boot () {
	cp /etc/wl/channellist /tmp
}
