# target/linux/atheros/wp838/wp838_QCA10.1.436/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS0::askfirst:/bin/ash --login
#ttyS0::askfirst:/usr/sbin/cli --login
tty0::askfirst:/bin/ash --login
#tty0::askfirst:/usr/sbin/cli --login
